# supervisor/nginx.conf
[program:nginx]
directory=/opt/nginx
command=nginx -c /opt/nginx/site.conf
exitcodes=0
autorestart=false
